# docker/images/php/container-files/etc/supervisor.d/php-fpm.conf
[program:php-fpm]
command = php-fpm
autorestart = true
stdout_logfile = /data/logs/php-fpm.log
stderr_logfile = /data/logs/php-fpm.log

[program:install-wordpress]
command = sh install-wp.sh
directory = /

[program:install-npm-dependiencies]
command = npm install
directory = /var/www