# docker/images/php/container-files/etc/supervisor.d/wordpress.conf
[program:install-wordpress]
command = sh install-wp.sh
directory = /

[program:install-npm-dependiencies]
command = npm install --unsafe-perm
directory = /var/www